# metac/metac.target
[Unit]
Description=MetaContainer
After=network.target
